# pgadmin4/rootfs/etc/services.d/pgadmin4/run
#!/usr/bin/with-contenv bashio

./entrypoint.sh &
pgadmin_pid = $!

# And let it roll
bashio::log.notice "pgAdmin 4 is running!"
wait "${pgadmin_pid}"
