# pgadmin4/rootfs/etc/services.d/pgadmin4/run
#!/usr/bin/with-contenv bashio

declare timeout

# Move to the right dir.
cd /pgadmin4

# Start Postfix to handle password resets etc.
sudo /usr/sbin/postfix start

# Get the session timeout from the pgAdmin config. We'll use this (in seconds)
# to define the Gunicorn worker timeout
timeout=$(cd /pgadmin4 && python -c 'import config; print(config.SESSION_EXPIRATION_TIME * 60 * 60 * 24)')

# Start pgAdmin
# NOTE: currently pgadmin can run only with 1 worker due to sessions implementation
# Using --threads to have multi-threaded single-process worker
exec gunicorn --timeout ${timeout} -w 1 --threads ${GUNICORN_THREADS:-25} --access-logfile ${GUNICORN_ACCESS_LOGFILE:--} run_pgadmin:app &

# Save the PID
pgadmin_pid=$!

# And let it roll
bashio::log.notice "pgAdmin 4 is running!"
wait "${pgadmin_pid}"
