# pgadmin4/rootfs/etc/services.d/pgadmin4/run
#!/usr/bin/with-contenv bashio

# Prepare SSL when needed
if bashio::config.true 'ssl'; then
    enable_tls=true
    mkdir -p /certs
    cp /ssl/$(bashio::config 'certfile') /certs/server.cert
    cp /ssl/$(bashio::config 'keyfile') /certs/server.key
fi


# Start pgAdmin
cd /pgadmin4
PGADMIN_ENABLE_TLS=$enable_tls \
    PGADMIN_DEFAULT_EMAIL="$(bashio::config 'default_email')" \
    PGADMIN_DEFAULT_PASSWORD="$(bashio::config 'default_password')" \
    /entrypoint.sh &

# Save the PID
pgadmin_pid=$!

# And let it roll
bashio::log.notice "pgAdmin 4 is running!"
wait "${pgadmin_pid}"
