# pgadmin4/rootfs/etc/services.d/pgadmin4/run
#!/usr/bin/with-contenv bashio

# Prepare SSL when needed
if bashio::config.true 'ssl'; then
    mkdir -p /certs
    cp /ssl/$(bashio::config 'certfile') /certs/server.cert
    cp /ssl/$(bashio::config 'keyfile') /certs/server.key
fi

# Start pgAdmin
PGADMIN_ENABLE_TLS=$(bashio::config 'ssl') \
    PGADMIN_DEFAULT_EMAIL="$(bashio::config 'default_email')" \
    PGADMIN_DEFAULT_PASSWORD="$(bashio::config 'default_password')" \
    ./entrypoint.sh &

# Save the PID
pgadmin_pid=$!

# And let it roll
bashio::log.notice "pgAdmin 4 is running!"
wait "${pgadmin_pid}"
