# pgadmin4/rootfs/etc/services.d/pgadmin4/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: pgAdmin4
# runs pgAdmin4 (see official-pgadmin4/pkg/docker/entrypoint.sh )
# ==============================================================================


declare timeout
declare pgadmin_pid

bashio::log.info "Starting pgAdmin.."

# Move to the right dir.
cd /pgadmin4

# Start Postfix to handle password resets etc.
sudo /usr/sbin/postfix start

# Get the session timeout from the pgAdmin config. We'll use this (in seconds)
# to define the Gunicorn worker timeout
timeout=$(cd /pgadmin4 && /venv/bin/python3 -c 'import config; print(config.SESSION_EXPIRATION_TIME * 60 * 60 * 24)')

# Start pgAdmin
# NOTE: currently pgadmin can run only with 1 worker due to sessions implementation
# Using --threads to have multi-threaded single-process worker
exec /venv/bin/gunicorn --timeout ${timeout} \
                        --bind unix:/tmp/pgadmin4.sock \
                        --workers=1 \
                        --threads=25 \
                        --chdir /pgadmin4 \
                        -c gunicorn_config.py \
                        pgAdmin4:app &
pgadmin_pid=$!

# Wait for pgAdmin4 to become available
bashio::log.info "Waiting for pgadmin to become available.."
while [ ! -S /tmp/pgadmin4.sock ]
do
    echo "."
    sleep 0.5
done

# And let it roll
bashio::log.notice "pgAdmin is running."
wait "${pgadmin_pid}"
